# launchd/com.willy.jarvis-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.willy.jarvis-backup</string>
    <key>Comment</key>
    <string>Nightly William Agent SQLite backup</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/willy/jarvis-core/scripts/backup-db.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/willy/jarvis-core</string>
    <key>StandardOutPath</key>
    <string>/Users/willy/jarvis-core/logs/backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/willy/jarvis-core/logs/backup.err.log</string>
  </dict>
</plist>
