# launchd/com.willy.jarvis-docker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.willy.jarvis-docker</string>
    <key>Comment</key>
    <string>William Agent Docker stack (Postgres, Redis, Caddy, backups)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/willy/jarvis-core/scripts/docker-stack-launchd.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/willy/jarvis-core</string>
    <key>StandardOutPath</key>
    <string>/Users/willy/jarvis-core/logs/docker-stack.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/willy/jarvis-core/logs/docker-stack.err.log</string>
  </dict>
</plist>
