# launchd/com.willy.william-kiosk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.willy.william-kiosk</string>
    <key>Comment</key>
    <string>William Agent fullscreen kiosk home screen</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/willy/jarvis-core/macos-helper/WilliamKiosk.app/Contents/MacOS/WilliamKiosk</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>AppleInterfaceStyle</key>
      <string>Dark</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/willy/jarvis-core/macos-helper</string>
    <key>StandardOutPath</key>
    <string>/Users/willy/jarvis-core/logs/kiosk.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/willy/jarvis-core/logs/kiosk.err.log</string>
  </dict>
</plist>
