# macos-helper/launchd/com.willy.jarvis-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.willy.jarvis-helper</string>
    <key>Comment</key>
    <string>William Agent macOS Helper</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/willy/jarvis-core/macos-helper/.build/release/JarvisHelper</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/willy/jarvis-core/macos-helper</string>
    <key>StandardOutPath</key>
    <string>/Users/willy/jarvis-core/logs/helper.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/willy/jarvis-core/logs/helper.err.log</string>
  </dict>
</plist>
