# macos-helper/launchd/com.willy.jarvis-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.willy.jarvis-helper</string>
    <key>Comment</key>
    <string>William Agent macOS Helper — MUST run signed JarvisHelper.app (not .build/release) for TCC</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/willy/jarvis-core/macos-helper/JarvisHelper.app/Contents/MacOS/JarvisHelper</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>WILLIAM_VOICE_BACKEND</key>
      <string>apple_legacy</string>
      <key>WILLIAM_OPENWAKEWORD_COMMAND</key>
      <string></string>
      <key>WILLIAM_OPENWAKEWORD_MODEL</key>
      <string></string>
      <key>WILLIAM_OPENWAKEWORD_THRESHOLD</key>
      <string>0.62</string>
      <key>WILLIAM_WHISPER_COMMAND</key>
      <string></string>
      <key>WILLIAM_WHISPER_MODEL</key>
      <string></string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/willy/jarvis-core/macos-helper</string>
    <key>StandardOutPath</key>
    <string>/Users/willy/jarvis-core/logs/helper.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/willy/jarvis-core/logs/helper.err.log</string>
  </dict>
</plist>
